# systemd/led-bike-lights.service
[Unit]
Description=LED Bike Wall Lights Controller
Documentation=file:///opt/led_bike_lights/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# NeoPixel (rpi_ws281x) needs root for DMA/PWM access.
User=root
Group=root
WorkingDirectory=/opt/led_bike_lights
Environment="PATH=/opt/led_bike_lights/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/led_bike_lights/venv/bin/python3 /opt/led_bike_lights/src/app.py
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5

# Hardening
PrivateTmp=yes
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-bike-lights

# Process limits
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target
